# opencontrail/files/3.0/contrail-vrouter-dpdk.ini
{%- from "opencontrail/map.jinja" import compute with context %}
[program:contrail-vrouter-dpdk]
command=taskset {{ compute.dpdk.taskset }} /usr/bin/contrail-vrouter-dpdk --no-daemon --vr_flow_entries=2097152 --vr_mpls_labels 6000 --socket-mem {{ compute.dpdk.socket_mem }}
priority=410
autostart=true
killasgroup=true
stdout_capture_maxbytes=1MB
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-vrouter-dpdk-stdout.log
stderr_logfile=/var/log/contrail/contrail-vrouter-dpdk-stderr.log
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
