# opencontrail/files/3.0/config/contrail-discovery

#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization API

supervisorctl -s unix:///tmp/supervisord_config.sock ${1} `basename ${0}:0`

#supervisorctl -s unix:///tmp/supervisord_config.sock ${1} `basename ${0}`
